# systemd/n8n-monitor.timer
[Unit]
Description=n8n Health Monitoring Timer
Requires=n8n-monitor.service

[Timer]
# Run every 5 minutes
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
